# packaging/systemd/coen-edge.service
[Unit]
Description=Coen edge (public tunnel endpoint)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/coen edge --config /etc/coen/edge.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2
User=coen
Group=coen
RuntimeDirectory=coen
RuntimeDirectoryMode=0750
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/coen /run/coen
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
